# deploy/systemd/tracecharter.service
[Unit]
Description=TraceCharter application server
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/tracecharter/current
Environment=NODE_ENV=production
Environment=PORT=8787
EnvironmentFile=-/etc/tracecharter/tracecharter.env
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=5
KillSignal=SIGINT
SyslogIdentifier=tracecharter

[Install]
WantedBy=multi-user.target
